# config/includes.chroot/etc/systemd/system/bw_ssh.service
[Unit]
Description=Enable ssh
[Service]
Type=simple
TimeoutStartSec=0
ExecStart=ufw allow ssh
[Install]
WantedBy=default.target
