# config/includes.chroot/etc/xdg/autostart/radeon.desktop
[Desktop Entry]
Version=1.0
Name=radeontop
Comment=Start radeontop in the background
Exec=radeontop -d /dev/null
Terminal=false
Type=Application
NoDisplay=true
